# configs/etc/cron.d/gitorama-com.production
MAILTO=svetlyak.40wt@gmail.com

# update stats every 15 minutes
*/15 * * * * art (cd /home/art/projects/gitorama.com && SETTINGS=gitorama.settings.production env/bin/python manage.py push_processes)
