# configs/etc/cron.d/gitorama-com.development
MAILTO=svetlyak.40wt@gmail.com

# update stats every hour
0 * * * * art (cd /vagrant/gitorama.com && ./manage.py update_users)

