# configs/etc/cron.d/gitorama-com.production
MAILTO=svetlyak.40wt@gmail.com

# update stats every hour
0 * * * * ubuntu (cd /home/ubuntu/projects/gitorama.com && ./manage.py update_users)

# update relations every day
0 4 * * * ubintu (cd /home/ubuntu/projects/gitorama.com && ./manage.py update_relations_stats)
